# setup/systemd/earlyoom.service.d/override.conf
# Install: sudo install -D -m 644 ~/hackery2/data/setup/systemd/earlyoom.service.d/override.conf \
#            /etc/systemd/system/earlyoom.service.d/override.conf
#          && sudo systemctl daemon-reload && sudo systemctl restart earlyoom
#
# The stock unit's DynamicUser=true implies RequiresMountsFor=/var/lib/private;
# on this machine var.mount is sometimes runtime-masked (seen 2026-07-02),
# which then blocks earlyoom from starting. Run as root instead (upstream's
# non-hardened default); CAP_KILL is inherent and the unit's other sandboxing
# (ProtectSystem=strict, ProtectHome) stays in effect.
[Service]
DynamicUser=no
