# data/authknock/authknock.service
[Unit]
Description=authknock

[Service]
ExecStart=/home/user/hackery2/data/authknock/main.py
Restart=always

[Install]
WantedBy=multi-user.target
