# setup/systemd/hdparm.timer
[Unit]
Description="hdparm"

[Timer]
AccuracySec=50s
OnBootSec=100
OnUnitActiveSec=150
Unit=hdparm.service

[Install]
WantedBy=multi-user.target
