# data/setup/machines/hp/etc/systemd/system/autonet.timer
[Unit]
Description="autonet"

[Timer]
OnBootSec=10
OnUnitActiveSec=10
Unit=autonet.service

[Install]
WantedBy=multi-user.target
