# data/setup/machines/hp/etc/systemd/system/pull_jj.service
[Unit]
Description="pull_jj"

[Service]
ExecStart=/home/koom/hackery2/src/hackery2/bin/pull_jj.sh

