# setup/machines/hp/etc/systemd/system/autonet.service
[Unit]
Description="autonet"

[Service]
ExecStart=/home/koom/hackery2/src/hackery2/bin/autonet.sh

