# setup/machines/hp/etc/systemd/system/obs_restarter.timer
[Unit]
Description="obs_restarter"

[Timer]
AccuracySec=0s
OnCalendar=*-*-* *:00:00
Unit=obs_restarter.service

[Install]
WantedBy=multi-user.target
