# setup/machines/hp/etc/systemd/system/obs_restarter.service
[Unit]
Description="obs_restarter"

[Service]
ExecStart=/home/koom/hackery2/src/hackery2/bin/obs_recording_restart.py
