# setup/systemd/plasma-xdg-desktop-portal-kde.service.d/override.conf
[Service]
# PyCharm (JetBrains Runtime) polls the Settings portal for GNOME keys every
# second; the Plasma 5.27 portal doesn't support that namespace and warns each
# time, flooding the journal. Discard exactly that message.
# Remove after upgrading to Plasma 6 (its portal serves org.gnome.desktop.interface).
# Deploy: ln -s /home/koom/hackery2/setup/systemd/plasma-xdg-desktop-portal-kde.service.d \
#   ~/.config/systemd/user/plasma-xdg-desktop-portal-kde.service.d
# Exact journald-side filter. NOT effective on systemd 255 (journald ignores
# LogFilterPatterns for user units there); kept so it activates after a systemd
# upgrade. Note: the whole value is the regex (no shell quoting), and the
# message has double spaces around the quoted namespace.
LogFilterPatterns=~^xdp-kde-settings: Namespace  "org[.]gnome[.]desktop[.]interface"  is not supported$

# Effective today: silence the warning at the source. Scoped to the portal's
# settings module logging category (the "xdp-kde-settings:" message prefix).
Environment=QT_LOGGING_RULES=xdp-kde-settings.warning=false
